# process-automated-builder/scripts/systemd/process-from-flow-batch.service
[Unit]
Description=Process-from-flow batch runner (persistent)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-%h/.openclaw/.env
EnvironmentFile=%h/.config/process-from-flow-batch/env
ExecStart=%h/.config/process-from-flow-batch/run-service.sh
Restart=always
RestartSec=5
KillMode=control-group
TimeoutStopSec=30

[Install]
WantedBy=default.target
